# scripts/crontab
0 1 * * * /opt/airmash/dehydrated/dehydrated -c >> /opt/airmash/logs/dehydrated.log 2>&1
